# .config/shell/profile.d/08-mpvopts.sh
#!/bin/sh

# mpv environment variables
export MPV_HOME="${XDG_CONFIG_HOME:-$HOME/.config}/mpv"
export MPV_LEAK_REPORT="1"
export DVDCSS_VERBOSE="1"
